# resources/debian/postinst
#!/bin/sh
# postinst script for rpistat

set -e

update-rc.d -f rpistat defaults

#DEBHELPER#

exit 0
